# L201utec.service
[Unit]
Description=Algo

[Service]
ExecStart=/usr/local/sbin/L201utec.sh

[Install]
WantedBy=multi-user.target
